# etc/local.d/bfq.start
#!/bin/bash
echo bfq > /sys/block/sda/queue/scheduler
